# deploy/systemd/music-library-worker.service
[Unit]
Description=Music Library Curator Celery worker
After=network.target postgresql.service redis-server.service
Requires=postgresql.service redis-server.service

[Service]
User=musiclibrary
Group=musiclibrary
WorkingDirectory=/opt/music-library-curator
EnvironmentFile=/opt/music-library-curator/.env
ExecStart=/opt/music-library-curator/.venv/bin/celery -A config worker --loglevel=INFO --concurrency=2
Restart=on-failure
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

